# files/etc/uci-defaults/01-setup-permission
#!/bin/sh
# Setup permission
opkg install /root/*.ipk --force-reinstall
chmod -R 755 /lib/*
chmod -R 755 /etc/profile.d/*
chmod 755 /etc/mutiara-wrt
chmod 755 /etc/profile
chown -R root:root .
