# files/etc/uci-defaults/01-setup-permission
#!/bin/sh

opkg install /root/*.ipk --force-reinstall

chmod -R 755 /lib/wifi/*

chmod -R 755 /lib/netifd/wireless/*

chown -R root:root .

chown -R mariadb:mariadb /srv/.

sudo service mysqld start
