# files/etc/uci-defaults/01-setup-permission
#!/bin/sh

opkg install /root/*.ipk --force-reinstall

chmod -R 755 /lib/*

chown -R root:root .
